# etc/cron.weekly/home-backup
#!/bin/bash
su jack -c "rdup-simple -E ~/.rdup-exclude ~ /media/SD/home-backup"
